# ansible/templates/home/pi/.config/systemd/user/set-wallpaper@.timer.j2
[Unit]
Description=Set the wallpaper after one minute.

[Timer]
OnBootSec=60s
Unit=set-wallpaper@%i.service

[Install]
WantedBy=timers.target
