# ansible/templates/etc/systemd/system/play-sound@.service.j2
[Unit]
Description=Play sound from /home/pi/Sounds

[Service]
Type=oneshot
ExecStart=/usr/bin/aplay -A 10 /home/pi/Sounds/%i.wav

[Install]
WantedBy=multi-user.target
