# ansible/templates/etc/systemd/system/timed-shutdown.timer.j2
[Unit]
Description=Timed shutdown at 17:05

[Timer]
OnCalendar=17:00:00

[Install]
WantedBy=timers.target
