# ansible/templates/etc/systemd/system/tv-off.service.j2
[Unit]
Description=Turn the TV off via CEC

[Service]
Type=oneshot
ExecStart=/usr/local/bin/remote-trigger tv_off
Environment=USER=root

[Install]
WantedBy=multi-user.target
